# ISBJava_Dizzy/trunk/tools/web/backupWiki.crontab
3 0 * * *        /users/sramsey/bin/backupWiki.pl 
