# ISBJava_Dizzy/trunk/tools/web/computeWebStats.crontab
MAILTO=sramsey
5 5 * * sun	/users/sramsey/bin/computeWebStats.pl Dizzy insDizzy ISBJava ISBJava\. Pointillist insPoint
